# deploy/sdg-ui.service
[Unit]
Description=SDG Digital Twin — Gradio UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/grad-project
EnvironmentFile=/opt/grad-project/.env
ExecStart=/opt/grad-project/.venv/bin/python3 /opt/grad-project/ui/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
